# tasks/crontab
# see: https://github.com/aptible/supercronic
# format: sec min hour dom month dow year command

# alerts: every 30 seconds
*/30 * * * * * * /usr/local/bin/alerts

# risk overview: every 30 minutes
0 */30 * * * * * /usr/local/bin/ghwo

# wpc probabilistic precip: hourly, 10 minutes past the hour (WaitForCycle
# still polls/retries internally in case WPC is running late)
0 10 * * * * * /usr/local/bin/wpcprob

# gridcache heat interval: every 30 minutes
0 */30 * * * * * /usr/local/bin/gridcache
